# DSM 6.X/Galaxy_NAT.conf
# /etc/init/Galaxy_NAT.conf
description "NAT with iptables"

author "Galaxy"

start on syno.network.ready

console log

script
	/usr/syno/etc/rc.sysv/Galaxy_NAT.sh start
end script

# vim:ft=upstart